# Artifact/arty/ingest_bot/deploy/backtest.timer
[Unit]
Description=Run ingest-bot backtest daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=backtest@btc.service

[Install]
WantedBy=timers.target
