# ingest_bot/deploy/ingest-bot.timer
[Unit]
Description=Run ingest-bot timer every hour

[Timer]
OnCalendar=hourly
Persistent=true
Unit=ingest-bot@prices.service

[Install]
WantedBy=timers.target
